# testgrounds/subjects/beta_evil/beta_evild.service

[Unit]
Description=beta_evild

[Service]
Type=simple
ExecStart=/usr/local/bin/beta_evild

[Install]
WantedBy=multi-user.target
