# redhat_bak/testgrounds/static/foo_daemon.service
[Unit]
Description=Foo's daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/foo_daemon

[Install]
WantedBy=multi-user.target
