# redhat_bak/testgrounds/static/happy_daemon.service

[Unit]
Description=Happy daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/happy_daemon

[Install]
WantedBy=multi-user.target
