# cwalld_test/cwalldtestd/cwalldtestd.service

[Unit]
Description=cwalldtestd

[Service]
Type=simple
ExecStart=/usr/local/bin/cwalldtestd

[Install]
WantedBy=multi-user.target
