# redhat_bak/testgrounds/testgrounds/static/bar_daemon.service
[Unit]
Description=Bar's daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/bar_daemon

[Install]
WantedBy=multi-user.target
