# cwalld/cwalld-enforce.service
[Unit]
Description=cwalld-enforce

[Service]
Type=simple
User=root
ExecStart=/usr/local/sbin/cwalld-enforce

[Install]
WantedBy=multi-user.target
