# testing/daemon_2.service

[Unit]
Description=Daemon2

[Service]
Type=simple
ExecStart=/usr/local/bin/daemon_2

[Install]
WantedBy=multi-user.target
