# cwalld_test/benchmark/cwalldspeedd/cwalldspeedd.service

[Unit]
Description=cwalldspeedd

[Service]
Type=simple
ExecStart=/usr/local/bin/cwalldspeedd

[Install]
WantedBy=multi-user.target
