# testgrounds/subjects/unconfined_daemon/unconfined_daemond.service

[Unit]
Description=unconfined_daemond

[Service]
Type=simple
ExecStart=/usr/local/bin/unconfined_daemond

[Install]
WantedBy=multi-user.target
