# testgrounds/subjects/unconfined_daemon/unconfined_daemon.service

[Unit]
Description=unconfined_daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/unconfined_daemon

[Install]
WantedBy=multi-user.target
