# testgrounds/subjects/gamma_civil/gamma_civild.service

[Unit]
Description=gamma_civild

[Service]
Type=simple
ExecStart=/usr/local/bin/gamma_civild

[Install]
WantedBy=multi-user.target
