# static_wall/declarations/daemon_4/daemon_4.service

[Unit]
Description=Daemon4

[Service]
Type=simple
ExecStart=/usr/local/bin/daemon_4

[Install]
WantedBy=multi-user.target
