# testgrounds/subjects/alpha_evil/alpha_evild.service

[Unit]
Description=alpha_evild

[Service]
Type=simple
ExecStart=/usr/local/bin/alpha_evild

[Install]
WantedBy=multi-user.target
