# testgrounds/subjects/gamma_confused/gamma_confusedd.service

[Unit]
Description=gamma_confusedd

[Service]
Type=simple
ExecStart=/usr/local/bin/gamma_confusedd

[Install]
WantedBy=multi-user.target
