# testgrounds/subjects/beta_civil/beta_civild.service

[Unit]
Description=beta_civild

[Service]
Type=simple
ExecStart=/usr/local/bin/beta_civild

[Install]
WantedBy=multi-user.target
