# testgrounds/subjects/alpha_curious/alpha_curiousd.service

[Unit]
Description=alpha_curiousd

[Service]
Type=simple
ExecStart=/usr/local/bin/alpha_curiousd

[Install]
WantedBy=multi-user.target
