# cwalld/cmd/cwalld/cwalld-enforce/cwalld-enforce.service

[Unit]
Description=cwalld-enforce

[Service]
Type=simple
ExecStart=/usr/local/bin/cwalld-enforce

[Install]
WantedBy=multi-user.target
